# rc.local
python3 /home/marshall/radio_button.py &
